# build_scripts/archive.cron
MAILTO="dev@code.org"

*/5 * * * * /etc/cdo/sync_dash
